# docker_templates/srv/howtox/supervisord.conf
[supervisord]
nodaemon=true
environment=HOME='/root'

[program:docker_start]
command=/bin/bash /srv/howtox/docker_start.sh
startsecs=0
autostart=true
autorestart=false

[program:c9]
command=/srv/howtox/cloud9/bin/cloud9.sh -l 0.0.0.0 -w /root/{{repoName}}

[program:web]
directory=/root/{{repoName}}
command=node /srv/howtox/web_server.js

[program:ttyterm]
directory=/root
command=/bin/bash

[program:tty]
directory=/srv/howtox
command=node ./tty_server.js
